# systemd/.config/systemd/user/git-fetch.timer
[Unit]
Description=Periodic git fetch for repos under ~/Work

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true
RandomizedDelaySec=30s

[Install]
WantedBy=timers.target
