# systemd/.config/systemd/user/git-fetch.service
[Unit]
Description=Background fetch all git repos under ~/Work
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/git-fetch-all
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
